# templates/themer.service.template
[Unit]
Description=Themer
[Service]
ExecStart=themer
Restart=always
[Install]
WantedBy=multi-user.target
